# business-rates-deb/src/deb/control/prerm
#!/bin/sh
if [ -x /bin/systemctl ]
then
    systemctl stop business-rates  || true
    systemctl disable business-rates  || true
else
    service business-rates stop || true
fi
