# business-rates-deb/src/deb/control/prerm
#!/bin/sh
systemctl stop business-rates || true
systemctl disable business-rates || true
